# releng/root-image/etc/systemd/system/startxfce4.service
[Unit]
Description=Start XFCE

[Service]
ExecStart=/usr/bin/startxfce4

[Install]
WantedBy=graphical.target