# releng/root-image/etc/rc.local
#!/bin/bash

chsh -s $(which zsh) root
chmod u+x /root/Desktop/install.sh